# example/shutdownWS_reboot.service
[Unit]
Description=restart
DefaultDependencies=no
RemainAfterExit=true

[Service]
ExecStart=/sbin/reboot
Type=oneshot