# example/shutdownWS_shutdown.service
[Unit]
Description=restart
DefaultDependencies=no
RemainAfterExit=true

[Service]
ExecStart=/sbin/shutdown -h now
Type=oneshot